# com.life-organizer.nightly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.life-organizer.nightly</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Update this path to your Python installation -->
        <string>/usr/bin/python3</string>
        <string>nightly.py</string>
        <string>--dry-run</string>
    </array>

    <key>WorkingDirectory</key>
    <!-- Update this path to your project directory -->
    <string>/Users/YOUR_USERNAME/Development/life-organizer</string>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- Run at 2:00 AM daily -->
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Development/life-organizer/nightly_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Development/life-organizer/nightly_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Add your API keys here if needed -->
        <!-- <key>OPENAI_API_KEY</key>
        <string>sk-...</string> -->
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
